# config/systemd/shoshop-selfsigned-ssl-renew.service
[Unit]
Description=Renew ShoShop self-signed TLS certificate if expiring soon
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/srv/shoshop
ExecStart=/srv/shoshop/scripts/ssl-selfsigned-renew.sh
# Путь замените, если проект не в /srv/shoshop
